# 17.Backup/borgbackup.timer
[Unit]
Description=BorgBackup
Requires=borgbackup.service

[Timer]
Unit=borgbackup.service
OnUnitActiveSec=300

[Install]
WantedBy=timers.target